# deploy/telegram-bot.service
[Unit]
Description=UstozHub notification queue processor (Telegram)
After=network.target redis-server.service postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/ustozhubuz
EnvironmentFile=/var/www/ustozhubuz/.env
ExecStart=/var/www/ustozhubuz/venv/bin/python manage.py process_notifications --daemon --interval 10
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
